# dot_config/systemd/user/flatpak-automatic.service
[Unit]
Description=Update Flatpak apps every 6 hours

[Service]
Type=simple
ExecStart=/usr/bin/flatpak update -y
